# deployment/configs/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:mysql]
command=/usr/bin/mysqld_safe
user=mysql
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/mysql.log

[program:backend]
command=java -jar /app/backend/target/*.jar
directory=/app
user=root
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/backend.log
environment=SPRING_PROFILES_ACTIVE="docker",SPRING_DATASOURCE_URL="jdbc:mysql://localhost:3306/search_assistant?useSSL=false&serverTimezone=Asia/Shanghai",SPRING_DATASOURCE_USERNAME="findthing",SPRING_DATASOURCE_PASSWORD="findthing123"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/nginx.log